# monitoring/systemd/boarderframeos-metrics.service
[Unit]
Description=BoarderframeOS Metrics Collector
After=network.target

[Service]
Type=simple
User=cosburn
WorkingDirectory=/Users/cosburn/BoarderframeOS
ExecStart=/opt/homebrew/opt/python@3.13/bin/python3.13 /Users/cosburn/BoarderframeOS/monitoring/metrics_collector.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
