# monitoring/systemd/boarderframeos-health.service
[Unit]
Description=BoarderframeOS Health Monitor
After=network.target

[Service]
Type=simple
User=cosburn
WorkingDirectory=/Users/cosburn/BoarderframeOS
ExecStart=/opt/homebrew/opt/python@3.13/bin/python3.13 /Users/cosburn/BoarderframeOS/monitoring/health_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
